# systemd/holo_mesh_sender.service
[Unit]
Description=HolographiX mesh sender
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /opt/holographix/examples/holo_mesh_sender.py --uri %i --chunk-dir /opt/holographix/data/%i.holo --peer 10.10.3.2:5000 --repeats 2
WorkingDirectory=/opt/holographix
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
